# srcs/influxdb/supervisor.conf
[supervisord]
nodaemon=true

[program:influxd]
command = sh -c "influxd && pkill supervisord "

[program:telegraf]
command = sh -c "telegraf --config /docker_telegraf.conf && pkill supervisord "