# srcs/nginx/srcs/config/supervisor.conf
[supervisord]
nodaemon=true

[program:sshd]
command = sh -c "/usr/sbin/sshd -D && pkill supervisord "

[program:nginx]
command = sh -c "nginx -g 'daemon off;' && pkill supervisord "